# deploy/banya-backend.service
[Unit]
Description=Banya FastAPI Backend
After=network.target postgresql.service

[Service]
User=root
WorkingDirectory=/root/banya/backend
Environment="PATH=/root/banya/backend/venv/bin:/usr/bin:/bin"
ExecStart=/root/banya/backend/venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
